# services.d/openvpn/run
#!/usr/bin/with-contenv bash

exec /usr/sbin/openvpn --script-security 2 --config /ovpn.conf --auth-user-pass /creds.txt